# recipes-support/gpioset-service/files/gpioset-low@.service
[Unit]
Description=Set GPIO low

[Service]
ExecStart=/usr/bin/gpioset %i=0
ExecStopPost=/usr/bin/gpioget %i

[Install]
WantedBy=multi-user.target

